# root/etc/s6-overlay/s6-rc.d/p11-firefox-customize/up
/etc/s6-overlay/s6-rc.d/p11-firefox-customize/run
